# Attacker-site/Attacker_site.conf
#!upstart
description "node.js server"
author      "David"

start on vagrant-mounted

script
	cd /vagrant
	exec sudo -u vagrant nodejs server.js
end script
