# MEAN-stack/MEAN_stack.conf
#!upstart
description "node.js server"
author      "Robin"

start on vagrant-mounted

script
	cd /vagrant
	exec sudo -u vagrant nodejs server.js
end script
